# docker/crontab
*/5 * * * * /app/monitor/src/monitor.py
