# docker/crontab
*/5 * * * * python /app/monitor/src/monitor.py
